# launchd/com.user.ollama-warm.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ollama-warm</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>for i in $(seq 1 30); do
  # Wait for the API to come up cheaply (no model load).
  if curl -fsS -m 5 http://127.0.0.1:11434/api/version >/dev/null 2>&amp;1; then
    break
  fi
  sleep 2
done
# Now do the actual warm. Model load can take 10-30s on a 16GB M4, but
# 2-3 minutes on smaller hosts where the model partially spills to CPU
# (e.g. 8GB M2 running the q4_K_M variant). Allow 300s.
if curl -fsS -m 300 -X POST http://127.0.0.1:11434/api/generate \
     -H 'Content-Type: application/json' \
     -d '{"model":"__MODEL_TAG__","prompt":"hi","stream":false,"keep_alive":-1}' \
     >/dev/null 2>&amp;1; then
  echo "warmed"
  exit 0
fi
echo "warmer failed"
exit 1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/ollama-warm.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/ollama-warm.err.log</string>
</dict>
</plist>
