# launchd/com.user.ollama-proxy.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ollama-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__REPO_ROOT__/proxy/ollama_proxy.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__BREW_PREFIX__/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>OLLAMA_PROXY_LISTEN_HOST</key>
        <string>0.0.0.0</string>
        <key>OLLAMA_PROXY_LISTEN_PORT</key>
        <string>11434</string>
        <key>OLLAMA_PROXY_UPSTREAM_HOST</key>
        <string>127.0.0.1</string>
        <key>OLLAMA_PROXY_UPSTREAM_PORT</key>
        <string>11435</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/ollama-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/ollama-proxy.err.log</string>
</dict>
</plist>
